# docker/root/etc/s6-overlay/s6-rc.d/svc-freshclam/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# freshclam daemon is disabled by default since Gamarr uses clamscan (standalone).
# Signatures are updated by init-clamav-config at startup.
# Set CLAMAV_DAEMON_ENABLED=true to run freshclam as a persistent daemon.
if [ "${CLAMAV_DAEMON_ENABLED,,}" != "true" ]; then
    # Update signatures once daily instead of running as daemon
    while true; do
        sleep 86400
        s6-setuidgid abc freshclam --quiet || true
    done
fi

exec \
    s6-setuidgid abc freshclam --daemon --foreground --checks=4
