# docker/root/etc/s6-overlay/s6-rc.d/svc-freshclam/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Update virus signatures daily. Runs as a lightweight loop rather than
# the freshclam daemon, since clamscan loads signatures on each scan.
# Initial database is downloaded by init-clamav-config at startup.
while true; do
    sleep 86400
    echo "**** updating clamav signatures ****"
    s6-setuidgid abc freshclam --quiet || echo "**** freshclam update failed ****"
done
