# docker/root/etc/s6-overlay/s6-rc.d/svc-clamd/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# clamd daemon disabled — Gamarr uses clamscan (standalone) which loads
# signatures on-demand per scan. No daemon needed.
# Set CLAMAV_DAEMON=true to run clamd if you need it for external tools.
if [ "${CLAMAV_DAEMON,,}" = "true" ]; then
    exec s6-setuidgid abc clamd --foreground
fi

exec sleep infinity
